# SupportFiles/LaunchAgents/com.github.execriez.labwarden.Sys-LoginWindowPoll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.execriez.labwarden.Sys-LoginWindowPoll</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/LabWarden/lib/Poll</string>
		<string>Sys-LoginWindow</string>
	</array>
	<key>LimitLoadToSessionType</key>
	<string>LoginWindow</string>
	<key>StartInterval</key>
	<integer>307</integer>
</dict>
</plist>
